# RabbitMQServer.conf
#!upstart
description "Rabbit MQ server"

start on startup
stop on shutdown

respawn

exec /home/krishan/git/project490/RabbitMQServer.php

post-start script
    PID=`status RabbitMQServer | egrep -oi '([0-9]+)$' | head -n1`
    echo $PID > /var/run/RabbitMQServer.pid
end script

post-stop script
    rm -f /var/run/RabbitMQServer.pid
end script

